# etc/service/socks-client/run
#!/bin/sh
export PORT=3013
cd /usr/app/client
exec yarn start:live