# etc/service/socks-server/run
#!/bin/sh
cd /usr/app/server
exec yarn start:live